# config/crontab
0 * * * *  bundle exec rake loomio:hourly_tasks
