# config/crontab
0 * * * * /loomio/bin/rake loomio:hourly_tasks > /proc/1/fd/1 2> /proc/1/fd/2
